# crontab
0 */3 * * * cd /var/www/html && php update-currencies.php >> /var/log/cron/cron.log 2>&1